# axi-test@.service
[Unit]
Description=Axi Test Instance (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/axi-tests/%i
ExecStart=/home/ubuntu/.local/bin/uv run python -m axi.supervisor
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
Environment="HOME=/home/ubuntu"
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
