# axi-bot.service
[Unit]
Description=Axi Personal Assistant Discord Bot
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
Environment="AXI_BOT_DIR=%h/axi-assistant"
WorkingDirectory=%E/AXI_BOT_DIR
ExecStart=/usr/bin/uv run python -m axi.supervisor
Restart=always
RestartSec=5
TimeoutStopSec=15
Environment="HOME=%h"
Environment="PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=default.target
